# alarmimg/usr/lib/systemd/system/oemcleanup.service
[Unit]
Description=Remove OEM user

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/oemcleanup

[Install]
WantedBy=multi-user.target
