# alarmimg/usr/lib/systemd/system/remove-calamares.service
[Unit]
Description=Remove Calamares

[Service]
Type=oneshot
ExecStart=/usr/bin/remove-calamares

[Install]
WantedBy=multi-user.target
